# systemd/belkinScan.timer
[Unit]
Description=Scan belkin daily

[Timer]
OnCalendar=daily
Persistent=false    
 
[Install]
WantedBy=timers.target
